# supervisord.conf
[supervisord]
nodaemon=true

[program:cron]
command=/usr/sbin/cron -f

[program:haproxy]
command=haproxy -f /etc/haproxy/haproxy.cfg
