# files/supervisord.conf
[supervisord]
nodaemon=true

[program:cron]
command=/usr/sbin/cron -f -L15

[program:haproxy]
command=/usr/local/bin/haproxy_run
